# deploy/systemd/pfs-api.service
[Unit]
Description=Mini Bloomberg — FastAPI Data Plane
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
WorkingDirectory=/opt/pfs
ExecStart=/root/.local/bin/uv run uvicorn src.api.app:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5
EnvironmentFile=/opt/pfs/.env

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/pfs/data

[Install]
WantedBy=multi-user.target
