# deploy/systemd/pfs-artifact-commit.timer
[Unit]
Description=Mini Bloomberg — Auto-commit artifacts to git

[Timer]
# Daily at 23:55 server local time
OnCalendar=*-*-* 23:55:00
Persistent=true

[Install]
WantedBy=timers.target
