# deploy/systemd/pfs-price-sync.timer
[Unit]
Description=Mini Bloomberg — Daily price sync (after market close)

[Timer]
# M-F 17:30 ET = 22:30 UTC (EST) or 21:30 UTC (EDT)
# Using America/New_York via OnCalendar with persistent
OnCalendar=Mon..Fri 21:30 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
