# deploy/systemd/pfs-streamlit.service
[Unit]
Description=Mini Bloomberg — Streamlit Presentation Plane
After=network.target pfs-api.service
Wants=pfs-api.service

[Service]
Type=simple
WorkingDirectory=/opt/pfs
ExecStart=/root/.local/bin/uv run streamlit run dashboard/app.py \
    --server.port 8501 \
    --server.address 100.106.13.112 \
    --server.headless true \
    --browser.gatherUsageStats false
Restart=on-failure
RestartSec=5
Environment=API_BASE_URL=http://127.0.0.1:8000

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/pfs/data /tmp

[Install]
WantedBy=multi-user.target
