# deploy/systemd/pfs-artifact-commit.service
[Unit]
Description=Mini Bloomberg — Auto-commit artifacts to git

[Service]
Type=oneshot
WorkingDirectory=/opt/pfs/data/artifacts
ExecStart=/opt/pfs/deploy/scripts/artifact-commit.sh

StandardOutput=journal
StandardError=journal
SyslogIdentifier=pfs-artifact-commit
