# deploy/systemd/pfs-mcp.service
[Unit]
Description=Mini Bloomberg - MCP HTTP Server
After=network-online.target pfs-api.service
Wants=network-online.target pfs-api.service

[Service]
Type=simple
WorkingDirectory=/opt/pfs
EnvironmentFile=/opt/pfs/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/root/.local/bin/uv run python -m pfs.mcp.server --http --host ${PFS_MCP_BIND_HOST} --port ${PFS_MCP_PORT}
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/pfs/data /tmp

[Install]
WantedBy=multi-user.target