# deploy/systemd/pfs-streamlit.service
[Unit]
Description=Mini Bloomberg — Streamlit Presentation Plane
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/pfs
EnvironmentFile=/opt/pfs/.env
ExecStart=/root/.local/bin/uv run streamlit run dashboard/app.py \
    --server.port ${PFS_STREAMLIT_PORT} \
    --server.address ${PFS_STREAMLIT_BIND_HOST} \
    --server.headless true \
    --browser.gatherUsageStats false
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/pfs/data /tmp

[Install]
WantedBy=multi-user.target
